# components/linux/v20260301/assets/check-route-overlap.service
[Unit]
Description=AKSFlexNode IPv4 route overlap pre-flight check
After=network-online.target static-routes.service
Wants=network-online.target
Before=kubelet.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /etc/aks-flex-node/check-route-overlap.sh
StandardOutput=journal
StandardError=journal

[Install]
RequiredBy=kubelet.service
