# components/linux/v20260301/assets/iptables-flush.service
[Unit]
Description=Flush iptables rules to a clean ACCEPT state
Before=kubelet.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/iptables-restore -cv /etc/aks-flex-node/iptables-clear
RemainAfterExit=no

[Install]
RequiredBy=kubelet.service
