# aks-flex-node-agent.service
[Unit]
Description=AKS Flex Node Agent
After=network-online.target
Wants=network-online.target
# Restart on failure to enable auto-recovery
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/usr/local/bin/aks-flex-node agent --config /etc/aks-flex-node/config.json
TimeoutStartSec=300
TimeoutStopSec=60
# Restart configuration for daemon resilience
Restart=on-failure
RestartSec=30
# AZURE_CONFIG_DIR points to a root-owned copy under /etc/aks-flex-node/azure,
# NOT the installing user's ~/.azure (which would allow local privilege escalation
# via malicious CLI extensions). See scripts/install.sh for the copy logic.
Environment=AZURE_CONFIG_DIR=/etc/aks-flex-node/azure
RuntimeDirectory=aks-flex-node
RuntimeDirectoryMode=0750
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
